# movie-api.service
[Unit]
Description=Movie System API Server
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/your/project
Environment=PATH=/usr/bin:/usr/local/bin
Environment=PYTHONPATH=/path/to/your/project
ExecStart=/usr/bin/python3 /path/to/your/project/api.py
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=movie-api

[Install]
WantedBy=multi-user.target
